# hosts/alpha/backup/pg-backup.timer
# /etc/systemd/system/pg-backup.timer
[Unit]
Description=Run pg-backup hourly

[Timer]
# Top-of-hour at :15 UTC, every hour. RandomizedDelaySec spreads the start
# to avoid colliding with other timers; Persistent catches missed runs
# (e.g. after the host was down) on next boot.
OnCalendar=*-*-* *:15:00 UTC
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
